# automation/cron/hh-gmail-monitor.cron
# Run HH Gmail monitor three times daily in local host time: 08:00, 14:00, 22:00.
# Replace <repo-root> with the absolute path when installing into crontab.
0 8,14,22 * * * cd <repo-root> && <repo-root>/automation/scripts/run-hh-gmail-monitor.sh
