# automation/systemd/personal-office-pi-job-search-gmail-monitor.service
[Unit]
Description=Personal Office Pi job-search Gmail monitor

[Service]
Type=oneshot
WorkingDirectory=/home/openclaw/personal-office-agent/personal-office
Environment=PERSONAL_OFFICE_ROOT=/home/openclaw/personal-office-agent/personal-office
Environment=OPENCLAW_JOB_SEARCH_TIMEOUT_SECONDS=1200
EnvironmentFile=-/home/openclaw/.config/personal-office/job-search-telegram.env
ExecStart=/home/openclaw/personal-office-agent/personal-office/automation/scripts/run-pi-job-search-gmail-monitor.sh
