# automation/systemd/personal-office-pi-archivist.timer
[Unit]
Description=Run Personal Office Pi archivist daily

[Timer]
OnCalendar=*-*-* 03:37:00
RandomizedDelaySec=20m
Persistent=true
Unit=personal-office-pi-archivist.service

[Install]
WantedBy=timers.target
