# automation/systemd/personal-office-hh-gmail-monitor.service
[Unit]
Description=Personal Office HH Gmail monitor

[Service]
Type=oneshot
WorkingDirectory=<repo-root>
ExecStart=<repo-root>/automation/scripts/run-hh-gmail-monitor.sh
